# conf.d/supervisor/slack-invite.conf
[program:slack_invite]
command=/go/bin/slack-invite
autostart=true
autorestart=true
stderr_logfile=/var/log/slack-invite.err.log
stdout_logfile=/var/log/slack-invite.out.log
